# doc/systemd/klangbecken-reload-jingles.timer
[Unit]
Description=Run klangbecken-reload-jingles.service daily

[Timer]
OnCalendar=*-*-* *:05:30

[Install]
WantedBy=timers.target
