# doc/systemd/klangbecken-fsck.timer
[Unit]
Description=Run klangbecken-fsck.service daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
