# doc/systemd/liquidsoap@klangbecken.service.d/overrides.conf
[Service]
EnvironmentFile=/etc/liquidsoap/klangbecken.env
IOSchedulingClass=best-effort
IOSchedulingPriority=1
CPUSchedulingPolicy=rr
CPUSchedulingPriority=1
