# doc/systemd/liquidsoap@klangbecken.service.d/overrides.conf
[Service]
EnvironmentFile=/etc/liquidsoap/klangbecken.env
IOSchedulingClass=best-effort
CPUSchedulingPolicy=rr
IOSchedulingPriority=1
CPUSchedulingPriority=90
